# etc/crontab
* * * * * ${HOME}/.local/src/cron/battery.sh
