# etc/crontab
* * * * * /home/josh/.local/src/cron/battery.sh
